# rpi-director-client.service
[Unit]
Description=Raspberry Pi LED Director Client
After=multi-user.target
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/boston/rpi-director
ExecStart=/home/boston/rpi-director-venv/bin/python3 -m rpi_director --mode client --client-id client1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Logging configuration - file logging disabled by default to protect SD card
Environment=RPI_DIRECTOR_ENABLE_FILE_LOGGING=0
Environment=RPI_DIRECTOR_LOG_LEVEL=WARNING

# To enable file logging, change RPI_DIRECTOR_ENABLE_FILE_LOGGING=1

[Install]
WantedBy=multi-user.target
